# systemd/pihole-mcp.service.example
# PiHole MCP Server - Systemd Service File
#
# This service runs the MCP server with HTTP transport on port 8765.
# It allows the analyzer.py agent to connect and make automated DNS analysis decisions.
#
# INSTALLATION:
# 1. Copy this file to /etc/systemd/system/pihole-mcp.service
# 2. Edit the paths below to match your installation
# 3. Run: sudo systemctl daemon-reload
# 4. Run: sudo systemctl enable --now pihole-mcp.service
#
# IMPORTANT: This assumes you have:
# - Python 3.11+ installed
# - uv package manager installed (https://github.com/astral-sh/uv)
# - PiHole running on the same system
# - The pihole-agent code in /home/pi5/pihole-agent (adjust paths as needed)

[Unit]
Description=PiHole MCP Server (HTTP Transport)
Documentation=https://modelcontextprotocol.io
After=network.target pihole-FTL.service
Wants=pihole-FTL.service

[Service]
Type=simple

# CHANGE THESE: Set to your user/group
User=pi5
Group=pi5

# CHANGE THIS: Set to your installation directory
WorkingDirectory=/home/pi5/pihole-agent

# Environment variables
# CHANGE THIS: Update HOME to your user's home directory
Environment=HOME=/home/pi5
Environment=PATH=/home/pi5/.local/bin:/usr/local/bin:/usr/bin:/bin

# The command to run the MCP server
# CHANGE THIS: Update paths if your installation location differs
ExecStart=/home/pi5/.local/bin/uv run python src/mcp_server/server.py

# Restart policy - automatically restart if it crashes
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only

# CHANGE THESE: Paths that need write access
# - logs: Where analysis reports are saved
# - data: Where any persistent data is stored (if needed)
# - .cache/uv: Required for uv package manager to function
# - .venv: Virtual environment (uv needs to update it)
ReadWritePaths=/home/pi5/pihole-agent/logs /home/pi5/pihole-agent/data /home/pi5/.cache/uv /home/pi5/pihole-agent/.venv

# Resource limits (optional, adjust as needed)
MemoryMax=256M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
